# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=npm start /opt/keeper
user=daemon
redirect_stderr=true
stdout_logfile=/var/opt/keeper/stdout.log
stdout_logfile_maxbytes=1MB
